# addons/servo_cam/rootfs/etc/services.d/servo-cam/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

export FLASK_HOST=0.0.0.0
export FLASK_PORT=$(bashio::config 'flask_port')

MODE=$(bashio::config 'mode')

if [[ "${MODE}" == "remote" ]]; then
    export REMOTE_HOST=$(bashio::config 'remote_host')
    export REMOTE_PORT=$(bashio::config 'remote_port')
    export REMOTE_SCHEME=$(bashio::config 'remote_scheme')

    if [[ -z "${REMOTE_HOST}" ]]; then
        bashio::log.error "Remote mode requires 'remote_host' to be configured."
        exit 1
    fi

    if [[ -z "${REMOTE_PORT}" || "${REMOTE_PORT}" == "null" ]]; then
        REMOTE_PORT=5000
    fi

    if [[ -z "${REMOTE_SCHEME}" || "${REMOTE_SCHEME}" == "null" ]]; then
        REMOTE_SCHEME="http"
    fi

    bashio::log.info "Starting Servo Cam remote proxy → ${REMOTE_SCHEME}://${REMOTE_HOST}:${REMOTE_PORT}"
    exec python3 /opt/servo_cam/remote_proxy.py
else
    if bashio::config.true 'flask_debug'; then
        export FLASK_DEBUG=true
    else
        export FLASK_DEBUG=false
    fi

    if bashio::config.has_value 'webhook_url'; then
        export WEBHOOK_URL=$(bashio::config 'webhook_url')
    fi

    bashio::log.info "Starting Servo Cam service on port ${FLASK_PORT}"
    exec python3 /opt/servo_cam/main.py
fi
